# examples/cron.d/currency
#
PATH=/home/stephens/local/src/currency/bin:/usr/local/bin:/usr/bin:/bin
MAILTO=void@asdf.com

0 12 * * *	nobody	currency_historical_rate_load -s xe 

